# lib/proxy/files/etc/crontabs/root
*/5 * * * * chmod +x /root/em/upgrade_check.sh ; /root/em/upgrade_check.sh > /root/em/log/upgrade_check.log